# deploy/com.afrog.btc-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.afrog.btc-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/marupro/CODEX/BTC_FX_CODEX/btc_monitor/.venv312/bin/python</string>
    <string>-u</string>
    <string>/Users/marupro/CODEX/BTC_FX_CODEX/btc_monitor/main.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/marupro/CODEX/BTC_FX_CODEX/btc_monitor</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/marupro/CODEX/BTC_FX_CODEX/btc_monitor/logs/runtime/monitor.out</string>

  <key>StandardErrorPath</key>
  <string>/Users/marupro/CODEX/BTC_FX_CODEX/btc_monitor/logs/runtime/monitor.err</string>
</dict>
</plist>
